# linuxOperation/supervisord/supervisord.conf.d/celery.ini
[group:linux-operation-celery]
programs=linux-operation-celery-worker, linux-operation-celery-beat

[program:linux-operation-celery-worker]
command                 = /usr/local/u-mail/app/engine/bin/celery -A operation worker -P gevent -c 4 -l info -Q linux-operation-celery-default
directory               = /usr/local/u-mail/operation
autostart               = true
autorestart             = true
redirect_stderr         = true
stdout_logfile          = /usr/local/u-mail/app/log/linux-operation-celery-worker.log
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups  = 2

[program:linux-operation-celery-beat]
command                 = /usr/local/u-mail/app/engine/bin/celery -A operation beat -l info
directory               = /usr/local/u-mail/operation
autostart               = true
autorestart             = true
redirect_stderr         = true
stdout_logfile          = /usr/local/u-mail/app/log/linux-operation-celery-beat.log
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups  = 2
